# docker-images/desktop-workspace/root/etc/s6-overlay/s6-rc.d/init-desktop-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Copy tint2 config if not exists
if [[ ! -f /config/.config/tint2/tint2rc ]]; then
    mkdir -p /config/.config/tint2
    cp /defaults/tint2rc /config/.config/tint2/tint2rc
fi

# Copy openbox config if not exists
if [[ ! -f /config/.config/openbox/menu.xml ]]; then
    mkdir -p /config/.config/openbox
    cp /defaults/menu.xml /config/.config/openbox/menu.xml
fi

if [[ ! -f /config/.config/openbox/autostart ]]; then
    mkdir -p /config/.config/openbox
    cp /defaults/autostart /config/.config/openbox/autostart
fi

# Configure GTK Bookmarks (File Manager Sidebar)
# Only add Datapool - File System (/) is shown by default in Thunar
if [[ ! -d /config/.config/gtk-3.0 ]]; then
    mkdir -p /config/.config/gtk-3.0
fi

BOOKMARKS_FILE="/config/.config/gtk-3.0/bookmarks"
# Create bookmarks with only Datapool (File System is built-in)
echo "file:///datapool Datapool" > "$BOOKMARKS_FILE"

# permissions
lsiown -R abc:abc \
    /config
lsiown -R abc:abc \
    /opt/obsidian
